# deploy/amazon-scraper.service
[Unit]
Description=AmazonScraper API Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/amazon-scraper
EnvironmentFile=/opt/amazon-scraper/.env
ExecStart=/opt/amazon-scraper/venv/bin/python api_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
